# units/simplehttpserver.service
[Unit]
Description=Python Simple HTTP Server
After=network.target

[Service]
ExecStart=/usr/bin/python -m SimpleHTTPServer 17291
Type=simple

[Install]
WantedBy=default.target
